# sampleapps/modified-images/mongodb/supervisord.conf
[supervisord]
nodaemon=true

[program:mongoexporter]
command=/bin/mongo_exporter

[program:mongo]
command=/entrypoint.sh mongod